# docker/app/overlay/etc/services.d/websockets/run
#!/bin/bash

exec chpst -u www-data:www-data php /var/www/html/artisan websockets:serve --debug
